# bin/ydns/ydns.timer
[Unit]
Description=Run ydns.io update every 15 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
