# bin/ydns/ydns.service
[Unit]
Description=Run ydns.io update

[Service]
Type=simple
ExecStart=/root/ydns_updater.sh

[Install]
WantedBy=default.target
